# runit/services/tty1/run
#!/bin/sh
exec chpst -P /sbin/agetty --autologin c --noclear 38400 tty1 linux
