# getty/getty@tty1.service.d/override.conf
[Service]
ExecStart=-/usr/bin/agetty --autologin cjoshi --noclear %I $TERM
