# runit/services/wlsunset/run
#!/bin/sh
exec 2>&1

exec wlsunset -l 33.8 -L -84.4
